# gameserver/templates/factorio.service.j2
[Unit]
Description={{ description }}

[Service]
ExecStart={{ exec_start }}
Type=simple
User={{ user }}
Group={{ user }}

[Install]
WantedBy=multi-user.target